# services/scripts/crontab
# Daily backup at 3 AM
0 3 * * * /backup.sh >> /backup/backup.log 2>&1
